# client/live-iso/config/includes.chroot/usr/lib/systemd/system/usbsas-rmmod.service
[Unit]
Description=rm USB modules for usbsas

[Service]
Type=oneshot
ExecStart=/usr/bin/usbsas-rmmod
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
